# roles/kodi/templates/kodi.service
[Unit]
Description=Kodi
After=remote-fs.target network-online.target
Wants=network-online.target

[Service]
User={{ kodi_user }}
Group={{ kodi_group }}
ExecStart=/usr/bin/kodi-standalone
Restart=always

[Install]
WantedBy=multi-user.target
